# rss-synchronizer.timer
[Unit]
Description=Regular rss-synchronizer update

[Timer]
OnCalendar=10m
